# v4/transform/MidrollByItem/MidrollByItem.cron
0 21 * * * python3 /root/v4/tansform/MidrollByItem/MidrollByItem.py >> /root/v4/tansform/MidrollByItem/cronlogs/filterMidrollByItem.log 2>&1
