# v4/transform/OverlayByItem/OverlayByItem.cron
0 21 * * * python3 /root/v4/tansform/OverlayByItem/OverlayByItem.py >> /root/v4/tansform/OverlayByItem/cronlogs/filterOverlayByItem.log 2>&1
