# v4/extract/PumpSlotStartLogs/PumpSlotStartLogs.cron
10 2 * * * python3 /root/v4/extract/PumpSlotStartLogs/PumpSlotStartLogs.py >> /root/v4/extract/PumpSlotStartLogs/cronlogs/PumpSlotStartLogs.log 2>&1
